# example3/docker/crontab/crontab
* * * * * /usr/bin/php /var/www/html/artisan schedule:run >> /dev/stdout 2>&1
